# supervisor.conf
[program:errbot]
command=INSERT_FULL_PATH_TO/run.sh
user=CHANGEME

autostart=true
autorestart=true
stderr_logfile=/var/log/errbot.err.log
stdout_logfile=/var/log/errbot.out.log
